# systemd/sift.service
[Unit]
Description=msgsift — fetch & classify mail into the daily todo store

[Service]
Type=oneshot
ExecStart=%h/.local/bin/msgsift
