# systemd/sift-nightly.timer
[Unit]
Description=Run msgsift nightly recap + rollover at end of day

[Timer]
OnCalendar=*-*-* 23:55:00
Persistent=true

[Install]
WantedBy=timers.target
