# systemd/sift.service
[Unit]
Description=msgsift — fetch & classify mail into the daily todo store
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/msgsift
