# systemd/sift-nightly.service
[Unit]
Description=msgsift — nightly recap + day rollover

[Service]
Type=oneshot
ExecStart=%h/.local/bin/msgsift-nightly
